# etc/cron.d/restic
SHELL=/bin/sh
PATH=/etc:/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin/:/usr/local/sbin/
# Order of crontab fields
# minute        hour    mday    month   wday    command
# Reference: https://www.freebsd.org/doc/handbook/configtuning-cron.html
# Reference: crontab(5).
@midnight	root	restic_backup.sh -c
@monthly	root	restic_check.sh
